# rpi2/etc/cron.d/hummingbird
# Collect reCamera temperature every 15 minutes
*/15 * * * * joe /opt/hummingbird/bin/collect_temp.sh
